# patch/lean/ezopwrt/rc.local
sleep 5 && wifi down
sleep 6 && wifi up
exit 0
