# patch/lean/ezopwrt/rc.local
sleep 5 && wifi down
sleep 10 && wifi up
exit 0
