# patch/lean/ezopwrt/rc.local
sleep 3 && wifi down
sleep 2 && wifi up
exit 0
